# puppet/modules/aws_ebs/templates/attach-volume.service.erb
[Unit]
Description=Attaches the EBS volume
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=<%= scope['aws_ebs::attach_bin_path'] %> /dev/<%= @device %> <%= @volume_id %>

[Install]
WantedBy=multiuser.target
