# puppet/modules/aws_ebs/templates/ensure-volume-formatted.service.erb
[Unit]
Description=Formats the ebs volume if needed
<% if @is_not_attached -%>
Requires=<%= @attach_service_name %>
After=<%= @attach_service_name %>
<% end -%>

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=<%= scope['aws_ebs::format_bin_path'] %> <%= @device %> <%= @filesystem %>

[Install]
WantedBy=multiuser.target
