# puppet/modules/aws_ebs/templates/volume.mount.erb
[Unit]
Description=Mount device to <%= @dest_path %>
Requires=<%= @format_service_name %>
After=<%= @format_service_name %>

[Mount]
What=<%= @device %>
Where=<%= @dest_path %>
Type=auto

[Install]
WantedBy=multiuser.target
