# puppet/modules/vault_client/templates/cert.timer.erb
[Unit]
Description=trigger certificate service <%= @name %>

[Timer]
OnUnitActiveSec=<%= @frequency %>
Persistent=true

[Install]
WantedBy=timers.target
