# puppet/modules/vault_client/templates/cert.service.erb
[Unit]
Description=certificate service <%= @name %>
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
Environment=VAULT_ADDR=http://127.0.0.1:8200
Environment=VAULT_DEV_ROOT_TOKEN_ID=root-token
Environment=VAULT_TOKEN=root-token
Environment=VAULT_CERT_TYPE=<%= @key_type %>
Environment=VAULT_CERT_BITS=<%= @key_bits %>
<% if @ip_sans.length > 0 -%>
Environment=VAULT_CERT_IP_SANS=<%= @ip_sans.join(',') %>
<% end -%>
ExecStart=/opt/bin/vault-helper cert \
  <%= @role %> \
  <%= @common_name %> \
  <%= @base_path %> \
<% if @alt_names.length > 0 -%>
  '--san-hosts=<%= @alt_names.join(',') %>' \
<% end -%>
<% if @ip_sans.length > 0 -%>
  '--ip-sans=<%= @ip_sans.join(',') %>' \
<% end -%>
  --owner=0 --group=0


<% @exec_post.each do |val| -%>
ExecStartPost=<%= val %>
<% end -%>

[Install]
WantedBy=multi-user.target
