# puppet/modules/vault_client/templates/secret.service.erb
[Unit]
Description=secret service <%= @name %>
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=/etc/vault/config
ExecStart=/opt/bin/vault-helper read <%= @secret_path %> -f <%= @field %> -d <%= @dest_path %> --owner=<%= @uid %> --group=<%= @gid %>
<% @exec_post.each do |val| -%>
ExecStartPost=<%= val %>
<% end -%>

[Install]
WantedBy=multi-user.target
