# puppet/modules/vault_client/templates/token-renewal.timer.erb
[Unit]
Description=trigger vault token renewal service

[Timer]
OnBootSec=60
OnUnitActiveSec=<%= @frequency %>
Persistent=true

[Install]
WantedBy=timers.target
