# puppet/modules/vault_client/templates/vault.service.erb
[Unit]
Description=Vault server
Requires=consul.service vault-assets.service
After=consul.service vault-assets.service

[Service]
Environment=AWS_REGION=<%= @region %>
ExecStart=/usr/local/bin/vault server \
-config=/etc/vault/vault.hcl
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
