# puppet/modules/kubernetes/templates/aws-iam-authenticator-init.service.erb
[Unit]
Description=Initialise heptio authenticator (AWS IAM authenticator)
Documentation=https://github.com/kubernetes-sigs/aws-iam-authenticator
After=network.target
Before=kube-apiserver.service
<%= scope.function_template(['kubernetes/_systemd_unit.erb']) %>
[Service]
Type=oneshot
WorkingDirectory=/var/aws-iam-authenticator
ExecStartPre=<%= scope['kubernetes::aia_bin_dir'] %>/aws-iam-authenticator \
  init -i <%= scope['kubernetes::cluster_name'] %>
ExecStartPre=/bin/mkdir -p /etc/kubernetes/aws-iam-authenticator
ExecStart=/bin/cp \
  /var/aws-iam-authenticator/aws-iam-authenticator.kubeconfig \
  <%= @auth_token_webhook_file %>

[Install]
WantedBy=multi-user.target
