# puppet/modules/kubernetes/templates/kube-apiserver.service.erb
[Unit]
Description=Kubernetes API Server
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
<%= scope.function_template(['kubernetes/_systemd_unit.erb']) %>

[Service]
Slice=podruntime.slice
User=<%= scope['kubernetes::user'] %>
Group=<%= scope['kubernetes::group'] %>
<%- if scope['kubernetes::_service_account_key_file'] and scope['kubernetes::service_account_key_generate'] -%>
ExecStartPre=/bin/sh -c "test -e <%= scope['kubernetes::_service_account_key_file'] %> || { openssl genrsa -out <%= scope['kubernetes::_service_account_key_file'] %> 2048 && chmod 640 <%= scope['kubernetes::_service_account_key_file'] %>; }"
<% end -%>
ExecStart=<%= scope['kubernetes::_dest_dir'] %>/apiserver \
  --v=<%= scope['kubernetes::log_level'] %> \
<% if @allow_privileged -%>
  --allow-privileged=true \
<% else -%>
  --allow-privileged=false \
<% end -%>
<% if @_audit_enabled -%>
  --audit-policy-file=<%= scope['kubernetes::apiserver::audit_policy_file'] %> \
  --audit-log-path=<%= scope['kubernetes::apiserver::audit_log_path'] %> \
  --audit-log-maxbackup=<%= scope['kubernetes::apiserver::audit_log_maxbackup'] %> \
  --audit-log-maxsize=<%= scope['kubernetes::apiserver::audit_log_maxsize'] %> \
<% end -%>
<% if @secure_port -%>
  --secure-port=<%= @secure_port %> \
<% end -%>
<% if @insecure_port -%>
  --insecure-port=<%= @insecure_port %> \
<% end -%>
<% if @insecure_bind_address -%>
  --insecure-bind-address=<%= @insecure_bind_address %> \
<% end -%>
<% if scope['kubernetes::cloud_provider'] != '' -%>
  --cloud-provider=<%= scope['kubernetes::cloud_provider'] %> \
<% end -%>
<% if scope['kubernetes::apiserver::ca_file'] -%>
  --client-ca-file=<%= scope['kubernetes::apiserver::ca_file'] %> \
<% end -%>
<% if scope['kubernetes::apiserver::cert_file'] -%>
  --tls-cert-file=<%= scope['kubernetes::apiserver::cert_file'] %> \
<% end -%>
<% if scope['kubernetes::apiserver::key_file'] -%>
  --tls-private-key-file=<%= scope['kubernetes::apiserver::key_file'] %> \
<% end -%>
<% if @kubelet_client_cert_file -%>
  --kubelet-client-certificate=<%= @kubelet_client_cert_file %> \
<% end -%>
<% if @kubelet_client_key_file -%>
  --kubelet-client-key=<%= @kubelet_client_key_file %> \
<% end -%>
<%- # add requestheaders + its CA for api aggregation support -%>
<%- if scope.function_versioncmp([scope['kubernetes::version'], '1.6.0']) >= 0 and \
        @requestheader_client_ca_file and \
        @proxy_client_cert_file and \
        @proxy_client_key_file -%>
  --requestheader-allowed-names=<%= @requestheader_allowed_names %> \
  --requestheader-extra-headers-prefix=<%= @requestheader_extra_headers_prefix %> \
  --requestheader-group-headers=<%= @requestheader_group_headers %> \
  --requestheader-username-headers=<%= @requestheader_username_headers %> \
  --requestheader-client-ca-file=<%= @requestheader_client_ca_file %> \
  --proxy-client-cert-file=<%= @proxy_client_cert_file %> \
  --proxy-client-key-file=<%= @proxy_client_key_file %> \
<%- end -%>
  "--etcd-servers=<%= @etcd_servers %>" \
  "--etcd-quorum-read=true" \
<% @etcd_servers_overrides.each do |override| -%>
  "--etcd-servers-overrides=<%= override %>" \
<% end -%>
<%- if scope['kubernetes::apiserver::etcd_ca_file'] -%>
  --etcd-cafile=<%= scope['kubernetes::apiserver::etcd_ca_file'] %> \
<% end -%>
<% if scope['kubernetes::apiserver::etcd_cert_file'] -%>
  --etcd-certfile=<%= scope['kubernetes::apiserver::etcd_cert_file'] %> \
<% end -%>
<% if scope['kubernetes::apiserver::etcd_key_file'] -%>
  --etcd-keyfile=<%= scope['kubernetes::apiserver::etcd_key_file'] %> \
<% end -%>
  --service-cluster-ip-range=<%= scope['kubernetes::service_ip_range_network'] %>/<%= scope['kubernetes::service_ip_range_mask'] %> \
  --admission-control=<%= @_admission_control.join(',') %> \
<%- if scope['kubernetes::_service_account_key_file'] -%>
  --service-account-key-file=<%= scope['kubernetes::_service_account_key_file'] %> \
<% end -%>
<%- if scope.function_versioncmp([scope['kubernetes::version'], '1.9.0']) >= 0 -%>
  --endpoint-reconciler-type=lease \
<%- else -%>
  --apiserver-count=<%= @count %> \
<%- end -%>
<% if @_storage_backend -%>
  --storage-backend=<%= @_storage_backend %> \
<% end -%>
<% if @authorization_mode.length > 0 -%>
  --authorization-mode=<%= @authorization_mode.join(',') %> \
<% end -%>
<% if @authorization_policy_file -%>
  --authorization-policy-file=<%= @authorization_policy_file %> \
<% end -%>
<% if @_runtime_config != [] -%>
  "--runtime-config=<%= @_runtime_config.join(',') %>" \
<% end -%>
<% if scope['kubernetes::apiserver::ca_file'] -%>
  "--kubelet-certificate-authority=<%= scope['kubernetes::apiserver::ca_file'] %>" \
<% end -%>
<% if @oidc_client_id -%>
  "--oidc-client-id=<%= @oidc_client_id %>" \
<% end -%>
<% if @oidc_groups_claim -%>
  "--oidc-groups-claim=<%= @oidc_groups_claim %>" \
<% end -%>
<% if @oidc_groups_prefix -%>
  "--oidc-groups-prefix=<%= @oidc_groups_prefix %>" \
<% end -%>
<% if @oidc_issuer_url -%>
  "--oidc-issuer-url=<%= @oidc_issuer_url %>" \
<% end -%>
<% if @_oidc_signing_algs != [] -%>
  "--oidc-signing-algs=<%= @_oidc_signing_algs.join(',') %>" \
<% end -%>
<% if @oidc_username_claim -%>
  "--oidc-username-claim=<%= @oidc_username_claim %>" \
<% end -%>
<% if @oidc_username_prefix -%>
  "--oidc-username-prefix=<%= @oidc_username_prefix %>" \
<% end -%>
  --logtostderr=true

Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
