# puppet/modules/kubernetes/templates/kubelet.service.erb
[Unit]
Description=Kubernetes Kubelet Server
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
<%= scope.function_template(['kubernetes/_systemd_unit.erb']) %>

[Service]
Slice=podruntime.slice
WorkingDirectory=<%= @kubelet_dir %>
<% if @cloud_provider == 'aws' -%>
# prevent metadata service access on AWS
ExecStartPre=/bin/sh -e -c "iptables -C PREROUTING -p tcp --destination 169.254.169.254 --dport 80 --in-interface <%= @container_interface %> --jump DNAT --table nat --to-destination 127.0.0.1:8181 2> /dev/null || iptables -A PREROUTING -p tcp --destination 169.254.169.254 --dport 80 --in-interface <%= @container_interface %> --jump DNAT --table nat --to-destination 127.0.0.1:8181"
<% end -%>
ExecStart=<%= scope['kubernetes::_dest_dir'] %>/kubelet \
  --v=<%= scope['kubernetes::log_level'] %> \
<% if scope.function_versioncmp([scope['kubernetes::version'], '1.6.0']) >= 0 -%>
<% if @_node_taints_string and @_node_taints_string.length > 0 -%>
  "--register-with-taints=<%= @_node_taints_string %>" \
<% end -%>
<% else -%>
  --register-schedulable=<%= @_register_schedulable %> \
<% end -%>
  --register-node=<%= @register_node %> \
<% if @allow_privileged -%>
  --allow-privileged=true \
<% else -%>
  --allow-privileged=false \
<% end -%>
<% if @kubeconfig_path -%>
  --kubeconfig=<%= @kubeconfig_path %> \
<% if scope.function_versioncmp([scope['kubernetes::version'], '1.4.0']) >= 0 -%>
<% if scope.function_versioncmp([scope['kubernetes::version'], '1.8.0']) < 0 -%>
  --require-kubeconfig \
<% end -%>
<% else -%>
  --api-servers=<%= scope['kubernetes::_master_url'] %> \
<% end -%>
<% end -%>
<% if @_node_labels_string and @_node_labels_string.length > 0 -%>
  "--node-labels=<%= @_node_labels_string %>" \
<% end -%>
<% if @hostname_override -%>
  --hostname-override=<%= @hostname_override %> \
<% end -%>
<% if @network_plugin -%>
  --network-plugin=<%= @network_plugin %> \
<% if @network_plugin_mtu -%>
  --network-plugin-mtu=<%= @network_plugin_mtu %> \
<% end -%>
<% end -%>
<% if @cloud_provider != '' -%>
  --cloud-provider=<%= @cloud_provider %> \
<% end -%>
<% if @container_runtime -%>
  --container-runtime=<%= @container_runtime %> \
<% end -%>
<% if not @kernelversion.nil? and scope.function_versioncmp([scope['kubernetes::version'], '1.6.0']) >= 0 and scope.function_versioncmp([scope['kubernetes::version'], '1.7.0']) < 0 and scope.function_versioncmp([@kernelversion, '4.9']) >= 0 -%>
  --cgroups-per-qos=false \
  --enforce-node-allocatable= \
<% end -%>
<% if @post_1_11 -%>
  --config=<%= @config_file  %> \
<% else -%>
  --cluster-dns=<%= @cluster_dns %> \
  --cluster-domain=<%= @cluster_domain %> \
<% if @pod_cidr -%>
  --pod-cidr=<%= @pod_cidr %> \
<% end -%>
<% if @client_ca_file and scope.function_versioncmp([scope['kubernetes::version'], '1.5.0']) >= 0 -%>
  --client-ca-file=<%= @client_ca_file %> \
  --anonymous-auth=false \
  --authentication-token-webhook \
  --authorization-mode=Webhook \
<% end -%>
<% # Kubernetes older than 1.7 has issues with kernel 4.9+ -%>
<% if scope.function_versioncmp([scope['kubernetes::version'], '1.6.0']) >= 0 -%>
  --cgroup-driver=<%= @cgroup_driver %> \
  --cgroup-root=<%= @cgroup_root %> \
<% if @cgroup_kube_name -%>
  --kube-reserved-cgroup=<%= @cgroup_kube_name %> \
  --runtime-cgroups=<%= @cgroup_kube_name %> \
  --kubelet-cgroups=<%= @cgroup_kube_name %> \
<%
    # build kube reserved command line
    @cgroup_kube_reserved = []
    @cgroup_kube_reserved << "cpu=#{@cgroup_kube_reserved_cpu}" unless @cgroup_kube_reserved_cpu.nil? or @cgroup_kube_reserved_cpu == 'nil'
    @cgroup_kube_reserved << "memory=#{@_cgroup_kube_reserved_memory}" unless @_cgroup_kube_reserved_memory.nil? or @_cgroup_kube_reserved_memory == 'nil'
    if @cgroup_kube_reserved.length > 0
-%>
  "--kube-reserved=<%= @cgroup_kube_reserved.join(',') %>" \
<% end -%>
<% end -%>
<% if @cgroup_system_name -%>
  --system-reserved-cgroup=<%= @cgroup_system_name %> \
<%
    # build system reserved command line
    @cgroup_system_reserved = []
    @cgroup_system_reserved << "cpu=#{@cgroup_system_reserved_cpu}" unless @cgroup_system_reserved_cpu.nil? or @cgroup_system_reserved_cpu == 'nil'
    @cgroup_system_reserved << "memory=#{@cgroup_system_reserved_memory}" unless @cgroup_system_reserved_memory.nil? or @cgroup_system_reserved_memory == 'nil'
    if @cgroup_system_reserved.length > 0
-%>
  "--system-reserved=<%= @cgroup_system_reserved.join(',') %>" \
<% end -%>
<% end -%>
<% end -%>
<% if @cert_file and @key_file -%>
  "--tls-cert-file=<%= @cert_file %>" \
  "--tls-private-key-file=<%= @key_file %>" \
<% end -%>
<%
    # build eviction hard command line
    @eviction_hard = []
    @eviction_hard << "memory.available<#{@eviction_hard_memory_available_threshold}" unless @eviction_hard_memory_available_threshold.nil? or @eviction_hard_memory_available_threshold == 'nil'
    @eviction_hard << "nodefs.available<#{@eviction_hard_nodefs_available_threshold}" unless @eviction_hard_nodefs_available_threshold.nil? or @eviction_hard_nodefs_available_threshold == 'nil'
    @eviction_hard << "nodefs.inodesFree<#{@eviction_hard_nodefs_inodes_free_threshold}" unless @eviction_hard_nodefs_inodes_free_threshold.nil? or @eviction_hard_nodefs_inodes_free_threshold == 'nil'
-%>
<% if @eviction_hard.length > 0 -%>
  "--eviction-hard=<%= @eviction_hard.join(',') %>" \
<% end -%>
<% if @eviction_soft_enabled -%>
<%
    # build eviction soft command line
    @eviction_soft = []

    @eviction_soft << "memory.available<#{@_eviction_soft_memory_available_threshold}" unless @_eviction_soft_memory_available_threshold.nil? or @_eviction_soft_memory_available_threshold == 'nil'
    @eviction_soft << "nodefs.available<#{@_eviction_soft_nodefs_available_threshold}" unless @_eviction_soft_nodefs_available_threshold.nil? or @_eviction_soft_nodefs_available_threshold == 'nil'
    @eviction_soft << "nodefs.inodesFree<#{@_eviction_soft_nodefs_inodes_free_threshold}" unless @_eviction_soft_nodefs_inodes_free_threshold.nil? or @_eviction_soft_nodefs_inodes_free_threshold == 'nil'
-%>
<%if @eviction_soft.length > 0 -%>
<%
    # build eviction soft grace period command line
    @eviction_soft_grace_period = []

    @eviction_soft_grace_period << "memory.available=#{@_eviction_soft_memory_available_grace_period}" unless @_eviction_soft_memory_available_grace_period.nil? or @_eviction_soft_memory_available_grace_period == 'nil'
    @eviction_soft_grace_period << "nodefs.available=#{@_eviction_soft_nodefs_available_grace_period}" unless @_eviction_soft_nodefs_available_grace_period.nil? or @_eviction_soft_nodefs_available_grace_period == 'nil'
    @eviction_soft_grace_period << "nodefs.inodesFree=#{@_eviction_soft_nodefs_inodes_free_grace_period}" unless @_eviction_soft_nodefs_inodes_free_grace_period.nil? or @_eviction_soft_nodefs_inodes_free_grace_period == 'nil'
-%>
<%        if @eviction_soft_grace_period.length > 0 -%>
  --eviction-soft=<%= @eviction_soft.join(',') %> \
  --eviction-soft-grace-period=<%= @eviction_soft_grace_period.join(',') %> \
  --eviction-max-pod-grace-period=<%= @eviction_max_pod_grace_period %> \
  --eviction-pressure-transition-period=<%= @eviction_pressure_transition_period %> \
<% end -%>
<% end -%>
<% end -%>
<%
    # build minumum reclaim command line
    @eviction_minimum_reclaim = []

    @eviction_minimum_reclaim << "memory.available=#{@eviction_minimum_reclaim_memory_available}" unless @eviction_minimum_reclaim_memory_available.nil? or @eviction_minimum_reclaim_memory_available == 'nil'
    @eviction_minimum_reclaim << "nodefs.available=#{@eviction_minimum_reclaim_nodefs_available}" unless @eviction_minimum_reclaim_nodefs_available.nil? or @eviction_minimum_reclaim_nodefs_available == 'nil'
    @eviction_minimum_reclaim << "nodefs.inodesFree=#{@eviction_minimum_reclaim_nodefs_inodes_free}" unless @eviction_minimum_reclaim_nodefs_inodes_free.nil? or @eviction_minimum_reclaim_nodefs_inodes_free == 'nil'
-%>
<% if @eviction_minimum_reclaim.length > 0 -%>
  "--eviction-minimum-reclaim=<%= @eviction_minimum_reclaim.join(',') %>" \
<% end -%>
<% if @_feature_gates != [] -%>
  --feature-gates=<%= @_feature_gates.join(',') %> \
<% end -%>
<% end -%>
  --logtostderr=true

Restart=on-failure
KillMode=process

[Install]
WantedBy=multi-user.target
