# puppet/modules/kubernetes/templates/kube-proxy.service.erb
[Unit]
Description=Kubernetes Kube-Proxy Server
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
After=network.target
<%= scope.function_template(['kubernetes/_systemd_unit.erb']) %>

[Service]
ExecStartPre=/sbin/sysctl -w net.bridge.bridge-nf-call-iptables=1
ExecStartPre=/sbin/sysctl -w net.bridge.bridge-nf-call-ip6tables=1
ExecStart=<%= scope['kubernetes::_dest_dir'] %>/proxy \
  --v=<%= scope['kubernetes::log_level'] %> \
  --cluster-cidr=<%= scope['kubernetes::pod_network'] %> \
<%
  if scope['kubernetes::kubelet::cgroup_kube_name']
  @cgroup_kube_basename = scope.call_function('regsubst', [scope['kubernetes::kubelet::cgroup_kube_name'], '^\/', ''])
-%>
  --resource-container=<%= @cgroup_kube_basename %> \
<%
  elsif scope['kubernetes::kubelet::cgroup_system_name']
  @cgroup_system_basename = scope.call_function('regsubst', [scope['kubernetes::kubelet::cgroup_system_name'], '^\/', ''])
-%>
  --resource-container=<%= @cgroup_system_basename %> \
<% end -%>
<% if @kubeconfig_path -%>
  --kubeconfig=<%= @kubeconfig_path %> \
<% end -%>
  --logtostderr=true

Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
