# puppet/modules/kubernetes/templates/kube-proxy.service.erb
[Unit]
Description=Kubernetes Kube-Proxy Server
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
After=network.target
<%= scope.function_template(['kubernetes/_systemd_unit.erb']) %>

[Service]
ExecStart=<%= scope['kubernetes::_dest_dir'] %>/proxy \
  --v=<%= scope['kubernetes::log_level'] %> \
  --cluster-cidr=<%= scope['kubernetes::pod_network'] %> \
  --resource-container=podruntime.slice \
<% if @kubeconfig_path -%>
  --kubeconfig=<%= @kubeconfig_path %> \
<% end -%>
  --logtostderr=true

Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
