# puppet/modules/kubernetes/templates/kubectl-apply.service.erb
[Unit]
Description=kubectl apply <%= @name %>
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
<%= scope.function_template(['kubernetes/_systemd_unit.erb']) %>

[Service]
User=<%= scope['kubernetes::user'] %>
Group=<%= scope['kubernetes::group'] %>
Type=oneshot
TimeoutStartSec=600
Environment=KUBECONFIG=<%= scope['kubernetes::kubectl::kubeconfig_path'] %>
RemainAfterExit=yes
# Wait for API server to be healthy
ExecStartPre=/bin/sh -c "while true; do STATUS=$(curl -s -o /dev/null -w '%{http_code}' http://localhost:8080/healthz); if [ $STATUS -eq 200 ]; then break; fi; done"
# Apply template
ExecStart=<%= @kubectl_path %> apply -f <%= @apply_file %>

[Install]
WantedBy=multi-user.target
