# puppet/modules/kubernetes/templates/kube-proxy.service.erb
[Unit]
Description=Kubernetes Kube-Proxy Server
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
After=network.target
<%= scope.function_template(['kubernetes/_systemd_unit.erb']) %>

[Service]
ExecStartPre=/sbin/sysctl -w net.bridge.bridge-nf-call-iptables=1
ExecStartPre=/sbin/sysctl -w net.bridge.bridge-nf-call-ip6tables=1
ExecStart=<%= scope['kubernetes::_dest_dir'] %>/proxy \
  --v=<%= scope['kubernetes::log_level'] %> \
<% if @post_1_11 -%>
  --config=<%= @config_file  %> \
<% else -%>
  --cluster-cidr=<%= scope['kubernetes::pod_network'] %> \
  --resource-container=podruntime.slice \
<% if @kubeconfig_path -%>
  --kubeconfig=<%= @kubeconfig_path %> \
<% end -%>
<% if @_feature_gates != [] -%>
  --feature-gates=<%= @_feature_gates.join(',') %> \
<% end -%>
  --logtostderr=true

Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
