# puppet/modules/kubernetes/templates/kube-scheduler.service.erb
[Unit]
Description=Kubernetes Scheduler Plugin
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
<%= scope.function_template(['kubernetes/_systemd_unit.erb']) %>

[Service]
Slice=podruntime.slice
User=<%= scope['kubernetes::user'] %>
Group=<%= scope['kubernetes::group'] %>
ExecStart=<%= scope['kubernetes::_dest_dir'] %>/scheduler \
  --v=<%= scope['kubernetes::log_level'] %> \
  --kubeconfig=<%= @kubeconfig_path %> \
  --leader-elect=true \
<% if @_feature_gates && @_feature_gates.length > 0 -%>
  --feature-gates=<% g = @_feature_gates.to_a.collect{|k| k.join('=')}.join(',') -%><%= g %> \
<% end -%>
  --profiling=false \
  --logtostderr=true

Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
