# puppet/modules/prometheus/files/localkube.service.d/netgo.conf
[Service]
Environment='GODEBUG=netdns=go'
