# puppet/modules/aws_es_proxy/templates/aws-es-proxy.service.erb
[Unit]
Description=AWS Elasticsearch Proxy
Documentation=https://github.com/abutaha/aws-es-proxy
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=<%= scope['proxy_path'] %> \
  -listen localhost:<%= scope['listen_port'] %> \
  -endpoint <%= scope['endpoint'] %>

Restart=on-failure
User=nobody
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
