# puppet/modules/consul/templates/consul.service.erb
[Unit]
Description=Consul server
<% if @_systemd_after.length > 0  -%>
After=<%= @_systemd_after.join(' ') %>
<% end  -%>
<% if @_systemd_before.length > 0  -%>
Before=<%= @_systemd_before.join(' ') %>
<% end  -%>
<% if @_systemd_wants.length > 0  -%>
Wants=<%= @_systemd_wants.join(' ') %>
<% end  -%>
<% if @_systemd_requires.length > 0  -%>
Requires=<%= @_systemd_requires.join(' ') %>
<% end  -%>

[Service]
User=<%= @user %>
Group=<%= @group %>
<% if @consul_bootstrap_expect > 1 -%>
Type=notify
TimeoutStartSec=600s
<% end -%>
PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p <%= @data_dir %>/data
ExecStartPre=/bin/chown -c root:<%= @group %> <%= @data_dir %>
ExecStartPre=/bin/chmod 0750 <%= @data_dir%>
ExecStartPre=/bin/chown -cR <%= @user %>:<%= @group %> <%= @data_dir %>/data
ExecStartPre=/bin/chmod 0750 <%= @data_dir %>/data
# detect an already existing datacenter setting
ExecStartPre=/bin/sh -c ". <%= @consul_detect_existing_datacenter_path %> <%= @data_dir %> > <%= @config_datacenter_override_path %>"
ExecStart=<%= @bin_path %> \
   agent \
   -config-file=<%= @config_path %> \
   -config-file=<%= @config_datacenter_override_path %> \
   -data-dir=<%= @data_dir %>/data
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
