# puppet/modules/consul/templates/consul-backup.timer.erb
[Unit]
Description=Backup consul according to schedule
[Timer]
OnCalendar=<%= @backup_schedule %>
[Install]
WantedBy=timers.target
