# puppet/modules/consul/templates/consul.service.erb
[Unit]
Description=Consul server
Requires=var-lib-consul.mount
After=var-lib-consul.mount

[Service]
User=consul
Group=consul
PermissionsStartOnly=true
ExecStartPre=/bin/mkdir -p /var/lib/consul/data
ExecStartPre=/bin/chown -c root:consul /var/lib/consul
ExecStartPre=/bin/chmod 0750 /var/lib/consul
ExecStartPre=/bin/chown -cR consul:consul /var/lib/consul/data
ExecStartPre=/bin/chmod 0750 /var/lib/consul/data
ExecStart=/usr/local/bin/consul \
   agent \
   -node=<%= @fqdn %> \
   -data-dir=/var/lib/consul/data \
   -server \
   -config-dir=/etc/consul \
   -advertise=<%= @private_ip %>

[Install]
WantedBy=multi-user.target
