# puppet/modules/consul/templates/var-lib-consul.mount.erb
[Unit]
Description=Mount consul data
After=ensure-ebs-volume-formatted.service
Requires=ensure-ebs-volume-formatted.service

[Mount]
What=/dev/xvdd
Where=/var/lib/consul
Type=ext4
