# puppet/modules/etcd/templates/etcd-backup.service.erb
[Unit]
Description=Backup etcd
Requires=etcd-k8s-events.service etcd-k8s-main.service etcd-overlay.service
After=etcd-k8s-events.service etcd-k8s-main.service etcd-overlay.service

[Service]
Type=oneshot
Environment=BUCKET_NAME=<%= @backup_bucket_prefix %>
Environment=CLUSTER_NAME=<%= @name %>
Environment=ENDPOINTS=<%= @endpoints %>
Environment=CA_PATH=<%= @ca_path %>
Environment=ETCD_VERSION=<%= @version %>
ExecStart=/opt/bin/etcd-backup.sh

[Install]
WantedBy=multi-user.target
