# puppet/modules/etcd/templates/etcd-backup.service.erb
[Unit]
Description=Backup etcd
Requires=etcd-k8s-events.service etcd-k8s-main.service etcd-overlay.service
After=etcd-k8s-events.service etcd-k8s-main.service etcd-overlay.service

[Service]
Type=oneshot
Environment=BUCKET_NAME=<%= @backup_bucket_prefix %>
Environment=OVERLAY_ENDPOINTS=<%= @overlay_endpoints %>
Environment=OVERLAY_CA=<%= @overlay_ca_name %>
Environment=K8S_EVENTS_ENDPOINTS=<%= @k8s_events_endpoints %>
Environment=K8S_EVENTS_CA=<%= @k8s_events_ca_name %>
Environment=K8S_MAIN_ENDPOINTS=<%= @k8s_main_endpoints %>
Environment=K8S_MAIN_CA=<%= @k8s_main_ca_name %>
ExecStart=/opt/bin/etcd-backup.sh

[Install]
WantedBy=multi-user.target
