# puppet/modules/vault_server/templates/vault-unsealer.service.erb
[Unit]
Description=Vault unsealer service
Requires=vault-assets.service
After=vault-assets.service

[Service]
Restart=always
RestartSec=5
Environment=VAULT_CACERT=/etc/vault/tls/ca.pem
Environment=VAULT_ADDR=https://localhost:8200
<% if @dev_mode -%>
Environment=VAULT_UNSEALER_MODE=local
Environment=VAULT_UNSEALER_LOCAL_KEY_DIR=<%= @vault_unsealer_key_dir %>
ExecStart=/opt/bin/vault-unsealer unseal
<% else -%>
Envirend=AWS_REGION=<%= @region %>
Environment=VAULT_UNSEALER_STORE_ROOT_TOKEN=false
Environment=VAULT_UNSEALER_MODE=aws-kms-ssm
Environment=VAULT_UNSEALER_AWS_KMS_KEY_ID=<%= @vault_unsealer_kms_key_id %>
Environment=VAULT_UNSEALER_AWS_SSM_KEY_PREFIX=<%= @vault_unsealer_ssm_key_prefix %>
ExecStart=/opt/bin/vault-unsealer unseal
<% end -%>

[Install]
WantedBy=multi-user.target
