# templates/blackbox_exporter.service.erb
[Unit]
Description=blackbox-exporter
After=network.target

[Service]
ExecStart=<%= @dest_dir %>/blackbox_exporter
Restart=on-failure

[Install]
WantedBy=multi-user.target
