# templates/calico-node.service.erb
[Unit]
Description=calico-node
<% if @systemd_after.length > 0  -%>
After=<%= @systemd_after.join(' ') %>
<% end  -%>
<% if @systemd_before.length > 0  -%>
Before=<%= @systemd_before.join(' ') %>
<% end  -%>
<% if @systemd_wants.length > 0  -%>
Wants=<%= @systemd_wants.join(' ') %>
<% end  -%>
<% if @systemd_requires.length > 0  -%>
Requires=<%= @systemd_requires.join(' ') %>
<% end  -%>

[Service]
EnvironmentFile=<%= scope['calico::params::config_dir'] %>/calico.env
ExecStartPre=-/usr/bin/docker rm -f calico-node
ExecStart=/usr/bin/docker run --net=host --privileged \
 --name=calico-node \
 -e HOSTNAME=${HOSTNAME} \
 -e IP=${CALICO_IP} \
 -e IP6=${CALICO_IP6} \
 -e CALICO_NETWORKING_BACKEND=${CALICO_NETWORKING_BACKEND} \
 -e AS=${CALICO_AS} \
 -e NO_DEFAULT_POOLS=${CALICO_NO_DEFAULT_POOLS} \
 -e CALICO_LIBNETWORK_ENABLED=${CALICO_LIBNETWORK_ENABLED} \
 -e ETCD_ENDPOINTS=${ETCD_ENDPOINTS} \
 -e ETCD_CA_CERT_FILE=${ETCD_CA_CERT_FILE} \
 -e ETCD_CERT_FILE=${ETCD_CERT_FILE} \
 -e ETCD_KEY_FILE=${ETCD_KEY_FILE} \
 -v /var/log/calico:/var/log/calico \
 -v /run/docker/plugins:/run/docker/plugins \
 -v /lib/modules:/lib/modules \
 -v /var/run/calico:/var/run/calico \
 -v <%= @etcd_cert_path %>:/etc/etcd/ssl \
 <%= scope['calico::params::calico_node_image'] %>:<%= @node_version%>

ExecStop=-/usr/bin/docker stop calico-node

[Install]
WantedBy=multi-user.target
