# templates/etcd-cert.timer.erb
[Unit]
Description=remint etcd-<%= @etcd_cluster %> certs daily and on boot

[Timer]
OnBootSec=1min
OnUnitActiveSec=1d
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
