# templates/kube-controller-manager.service.erb
[Unit]
Description=Kubernetes Controller Manager
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
<%= scope.function_template(['kubernetes/_systemd_unit.erb']) %>

[Service]
User=<%= scope['kubernetes::user'] %>
Group=<%= scope['kubernetes::group'] %>
<%- if scope['kubernetes::_service_account_key_file'] and scope['kubernetes::service_account_key_generate'] -%>
ExecStartPre=/bin/sh -c "test -e <%= scope['kubernetes::_service_account_key_file'] %> || { openssl genrsa -out <%= scope['kubernetes::_service_account_key_file'] %> 2048 && chmod 640 <%= scope['kubernetes::_service_account_key_file'] %>; }"
<% end -%>
ExecStart=<%= scope['kubernetes::_dest_dir'] %>/controller-manager \
  --v=<%= scope['kubernetes::log_level'] %> \
<% if scope['kubernetes::cloud_provider'] != '' -%>
  --cloud-provider=<%= scope['kubernetes::cloud_provider'] %> \
<% end -%>
  --kubeconfig=<%= @kubeconfig_path %> \
<% if scope['kubernetes::controller_manager::ca_file'] -%>
  --root-ca-file=<%= scope['kubernetes::controller_manager::ca_file'] %> \
<% end -%>
<%- if scope['kubernetes::_service_account_key_file'] -%>
  --service-account-private-key-file=<%= scope['kubernetes::_service_account_key_file'] %> \
<% end -%>
  --leader-elect=true \
  --logtostderr=true

Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
