# templates/kube-scheduler.service.erb
[Unit]
Description=Kubernetes Scheduler Plugin
Documentation=https://github.com/GoogleCloudPlatform/kubernetes

[Service]
User=<%= scope['kubernetes::user'] %>
Group=<%= scope['kubernetes::group'] %>
ExecStart=<%= scope['kubernetes::real_dest_dir'] %>/controller-manager \
  --v=<%= scope['kubernetes::log_level'] %> \
  --kubeconfig=<%= @kubeconfig_path %> \
  --logtostderr=true

Restart=on-failure
Type=notify
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
