# templates/kubelet.service.erb
[Unit]
Description=Kubernetes Kubelet Server
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
<% if $container_engine == 'docker' %>
After=docker.service
Requires=docker.service
<% end -%>

[Service]
WorkingDirectory=<%= @kubelet_dir %>
ExecStart=<%= scope['kubernetes::_dest_dir'] %>/kubelet \
  --v=<%= scope['kubernetes::log_level'] %> \
  --register-schedulable=<%= @_register_schedulable %> \
  --register-node=<%= @register_node %> \
<% if @allow_privileged -%>
  --allow-privileged=true \
<% else -%>
  --allow-privileged=false \
<% end -%>
<% if @kubeconfig_path -%>
  --kubeconfig=<%= @kubeconfig_path %> \
<% if scope.function_versioncmp([scope['kubernetes::version'], '1.4.0']) >= 0 -%>
  --require-kubeconfig \
<% else -%>
  --api-servers=<%= scope['kubernetes::master_url'] %> \
<% end -%>
<% end -%>
<% if @_node_labels_string and @_node_labels_string.length > 0 -%>
  --node-labels=<%= @_node_labels_string %> \
<% end -%>
  --cluster-dns=<%= @cluster_dns %> \
  --cluster-domain=<%= @cluster_domain %> \
<% if @hostname_override -%>
  --hostname-override=<%= @hostname_override %> \
<% end -%>
<% if @pod_cidr -%>
  --pod-cidr=<%= @pod_cidr %> \
<% end -%>
<% if @network_plugin -%>
  --network-plugin=<%= @network_plugin %> \
<% if @network_plugin_mtu -%>
  --network-plugin-mtu=<%= @network_plugin_mtu %> \
<% end -%>
<% end -%>
<% if scope['kubernetes::cloud_provider'] != '' -%>
  --cloud-provider=<%= scope['kubernetes::cloud_provider'] %> \
<% end -%>
<% if @container_runtime -%>
  --container-runtime=<%= @container_runtime %> \
<% end -%>
<% if @_ca_file and scope.function_versioncmp([scope['kubernetes::version'], '1.5.0']) >= 0 -%>
  --client-ca-file=<%= @_ca_file %> \
<% end -%>
<% # Kubernetes older than 1.7 has issues with kernel 4.9+ -%>
<% if not @kernelversion.nil? and scope.function_versioncmp([scope['kubernetes::version'], '1.6.0']) >= 0 and scope.function_versioncmp([scope['kubernetes::version'], '1.7.0']) < 0 and scope.function_versioncmp([@kernelversion, '4.9']) >= 0 -%>
  --cgroups-per-qos=false \
  --enforce-node-allocatable= \
<% end -%>
<% if scope.function_versioncmp([scope['kubernetes::version'], '1.6.0']) >= 0 -%>
  --cgroup-driver=<%= @cgroup_driver %> \
<% end -%>
  --logtostderr=true

Restart=on-failure
KillMode=process

[Install]
WantedBy=multi-user.target
