# templates/kubectl-apply.service.erb
[Unit]
Description=kubectl apply <%= @name %>
Documentation=https://github.com/GoogleCloudPlatform/kubernetes
<%= scope.function_template(['kubernetes/_systemd_unit.erb']) %>

[Service]
User=<%= scope['kubernetes::user'] %>
Group=<%= scope['kubernetes::group'] %>
Type=oneshot
TimeoutStartSec=600
# Wait for API server to be healthy
ExecStartPre=<%= @curl_path %> --retry 120 --retry-delay 5 -s  http://localhost:8080/healthz
# Apply template
ExecStart=<%= @kubectl_path %> apply -f <%= @apply_file %>

[Install]
WantedBy=multi-user.target
