# templates/prometheus-node-exporter.service.erb
[Unit]
Description=prometheus-node-exporter
After=docker.service
Requires=docker.service

[Service]
ExecStartPre=-/usr/bin/docker rm -f prometheus-node-exporter
ExecStart=/usr/bin/docker run \
 --name=prometheus-node-exporter \
 -p <%= @node_exporter_port %>:<%= @node_exporter_port %> \
 -v /proc:/host/proc:ro \
 -v /dev:/host/dev:ro \
 -v /sys:/host/sys:ro \
 -v /:/rootfs:ro \
 <%= @node_exporter_image %>:<%= @node_exporter_version %> \
 -collector.procfs \
 /host/proc \
 -collector.sysfs \
 /host/sys \
 -collector.filesystem.ignored-mount-points \
 '"^/(sys|proc|dev|host|etc)($|/)"'

ExecStop=-/usr/bin/docker stop prometheus-node-exporter

[Install]
WantedBy=multi-user.target
