# templates/secret.service.erb
[Unit]
Description=secret service <%= @name %>
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=VAULT_CERT_OWNER=<%= @user %>:<%= @group %>
ExecStart=/etc/vault/helper read <%= @secret_path %> <%= @field %> <%= @dest_path %>
<% @exec_post.each do |val| -%>
ExecStartPost=<%= val %>
<% end -%>

[Install]
WantedBy=multi-user.target
