# templates/secret.service.erb
[Unit]
Description=secret service <%= @name %>
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=VAULT_CERT_OWNER=<%= @user %>:<%= @group %>
Environment=VAULT_ADDR=http://127.0.0.1:8200
Environment=VAULT_DEV_ROOT_TOKEN_ID=init-token-client
Environment=VAULT_TOKEN=init-token-client
ExecStart=/opt/bin/vault-helper read <%= @secret_path %> -f <%= @field %> -d <%= @dest_path %> --owner=0 --group=0
<% @exec_post.each do |val| -%>
ExecStartPost=<%= val %>
<% end -%>

[Install]
WantedBy=multi-user.target
