# templates/token-renewal.service.erb
[Unit]
Description=vault token renewal service
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
Environment=VAULT_ADDR=http://127.0.0.1:8200
Environment=VAULT_DEV_ROOT_TOKEN_ID=root-token
Environment=VAULT_TOKEN=root-token
Environment=VAULT_CLUSTER_ID=test
ExecStart=/opt/bin/vault-helper renew-token $VAULT_CLUSTER_ID --role=master

[Install]
WantedBy=multi-user.target
