# templates/token-renewal.service.erb
[Unit]
Description=vault token renewal service
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/vault/config
ExecStart=/opt/bin/vault-helper renew-token --init_role=$(VAULT_INIT_ROLE)

[Install]
WantedBy=multi-user.target
