# templates/token-renewal.timer.erb
[Unit]
Description=trigger vault token renewal service

[Timer]
OnUnitActiveSec=<%= @frequency %>
Persistent=true

[Install]
WantedBy=timers.target
