# templates/token-renewal.service.erb
[Unit]
Description=vault token renewal service
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
ExecStart=/opt/vault-helper renew-token <%= @ec2_tag_cluster %>

[Install]
WantedBy=multi-user.target
