# templates/token-renewal.service.erb
[Unit]
Description=vault token renewal service
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=oneshot
ExecStart=/etc/vault/helper token-renew
